# 99-usbbackup.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sd[a-z]1", ATTRS{idProduct}=="0bc2", ATTRS{idVendor}=="3322", RUN+="/usr/local/bin/usbbackup"
